# ignition/files/nazar-rpm-layer.service
[Unit]
Description=Nazar First-Boot RPM Layer
ConditionPathExists=!/etc/nazar/.rpm-layer-complete
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rpm-ostree install --idempotent --allow-inactive yq jq git-core nodejs22 tailscale vim-minimal htop tmux
ExecStartPost=/usr/bin/touch /etc/nazar/.rpm-layer-complete
ExecStartPost=/usr/bin/systemctl reboot
RemainAfterExit=yes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
