# image/root/etc/systemd/system/nazar-evolve-resume.service
[Unit]
Description=Nazar Evolution Resume
After=network-online.target nazar-setup.service
ConditionPathExists=/var/lib/nazar/evolution/pending.yaml

[Service]
Type=oneshot
ExecStart=/usr/bin/nazar evolve --resume
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
